# chromium/s6-overlay/s6-rc.d/cdp-proxy/run
#!/command/with-contenv bash
# s6-Service: CDP-Proxy
exec python3 /usr/local/bin/cdp-proxy.py
